# deploy/linux/ao2-cp-server.service
[Unit]
Description=ao2-control-plane server (read-only observer for AO2 evidence)
Documentation=https://github.com/uesugitorachiyo/ao2-control-plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ao2cp
Group=ao2cp
WorkingDirectory=/var/lib/ao2-cp-server
ExecStart=/usr/local/bin/ao2-cp-server --bind 127.0.0.1:8744 --data-dir /var/lib/ao2-cp-server/data
EnvironmentFile=/etc/ao2-cp-server.env
Restart=on-failure
RestartSec=5

# Process hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

# Writable paths (data dir + optional release-smoke artifact dir)
ReadWritePaths=/var/lib/ao2-cp-server/data

# Network restrictions: server binds 127.0.0.1 only; deny IPv6/raw sockets
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# File-descriptor headroom for 500-concurrent burst tests / load
# (matches the ulimit -n 65536 used in CI). Server steady-state is far lower.
LimitNOFILE=65536

# Forbid the server from inheriting any of the env vars the preflight rejects.
# (AO2_CP_API_TOKEN comes from EnvironmentFile; OPENAI_API_KEY / ANTHROPIC_API_KEY
# must NOT be present or the server will refuse to start.)
PassEnvironment=

[Install]
WantedBy=multi-user.target
